# roles/sentry/templates/supervisor-sentry.conf.j2
[program:sentry]
command=/var/sentry/ve/bin/sentry --config=/var/sentry/sentry_conf.py start http
environment=PATH="/var/sentry/ve/bin"
autostart=true
autorestart=true
redirect_stderr=true
